# deploy/systemd/weather-map-poll-tiles.service
[Unit]
Description=Poll tiles from artifacts and restart Martin if data changes
Wants=network-online.target cloud-final.service
After=network-online.target cloud-final.service
ConditionPathExists=/etc/weather-map/poll-tiles.env

[Service]
Type=oneshot
EnvironmentFile=-/etc/weather-map/poll-tiles.env
ExecStart=/opt/weather-map/scripts/poll-tiles.sh
